# supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisor.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=0            ; (num of main logfile rotation backups;default 10)
loglevel=warn                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:Apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
autorestart=true
priority=1

[program:voipmonitor]
command=/usr/local/sbin/voipmonitor --config-file=/etc/voipmonitor.conf --pid-file=/var/run/voipmonitorinitscript.pid -v 1
autorestart=false
priority=2